# ops/crontab
# 容器内系统 cron 调度表（supercronic 解析；它继承容器环境变量并传给任务）。
# 字段：分 时 日 月 周  命令。触发脚本用 Node fetch 打 /api/cron（免 curl）。
# 真正的采集→分析→校验→brief 在 app 容器进程内跑（非 serverless，无超时），
# 由 /api/cron 走 Job Runner 落 Run / 成本 / 审计。

# 每日 1 次完整管线（dogfood feedback 2026-06-06：spec 写"Daily Brief"但之前是
# 每 6h 一份，一天 4 份 brief 名不副实；改为每日凌晨 1:00 北京 = UTC 17:00 跑，
# 早晨用户可读到当日 brief。源 fetch_interval=6h 是采集层"想多久抓一次"的提示，
# 与 cron 触发频率无关——前者描述意愿、后者真触发；不一致已在 source_id
# fetch_interval 字段加文档注解）
0 17 * * * node --no-warnings /app/ops/trigger.mjs

# 每日卷备份：管线（17:00 UTC）跑完后 1 小时备份，确保当日 brief 已落库再快照。
# SQLite 在线备份（一致）+ reports/ 正文 → /data/backups/<时间戳>/，默认留最近 14 份
# （BACKUP_KEEP 可调）。备份在同卷——防损坏/误删，非防整卷丢失；真 DR 把 /data/backups
# 拉到机器外（见 operations.md §6）。需 cron 服务挂 /data 卷（docker-compose.yml 已配）。
0 18 * * * node --no-warnings /app/ops/backup-db.mjs
